# 999-xmodmap.rules
ACTION=="add", RUN+="/home/icook/bin/remap-kbd"
